# systemd/go2-vercel.service
[Unit]
Description=Go2 camera, LiDAR and teleoperation bridge to Vercel/LiveKit
Wants=network-online.target
After=NetworkManager.service network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=unitree
Group=unitree
WorkingDirectory=/home/unitree/Desktop/SLAM-Navigation---GO2-EDU-Mapping-Rviz-and-WebViewr-Map-
Environment=HOME=/home/unitree
Environment=PATH=/home/unitree/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=OPEN_BROWSER=0
Environment=PYTHONUNBUFFERED=1
Environment=GO2_CAMERA_INTERFACE=eth0
Environment="GO2_UPLINK_CONNECTION=XD4 Local"
Environment=GO2_UPLINK_SSID=XD4
Environment=GO2_FALLBACK_CONNECTION=Thaina
Environment=GO2_FALLBACK_SSID=Thaina
Environment=GO2_UPLINK_INTERFACE=wlan0
Environment=GO2_ROBOT_INTERFACE=eth0
EnvironmentFile=-/etc/oracle-go2-teleoperation/network.env

# Prefere o SSID XD4, usa Thaina como fallback e mantem a Ethernet do Go2.
ExecStartPre=+/usr/bin/bash /home/unitree/Desktop/SLAM-Navigation---GO2-EDU-Mapping-Rviz-and-WebViewr-Map-/scripts/prepare_thaina_network.sh
ExecStart=/usr/bin/bash /home/unitree/Desktop/SLAM-Navigation---GO2-EDU-Mapping-Rviz-and-WebViewr-Map-/run_vercel.sh

# Se a internet ainda nao estiver pronta, run_vercel.sh encerra e o systemd tenta
# novamente. Isso tambem recupera quedas do LiveKit sem intervencao manual.
Restart=always
RestartSec=10
TimeoutStopSec=25
KillMode=control-group

[Install]
WantedBy=multi-user.target
